# docker/php/supervisord.conf
[supervisord]
nodaemon=true

[program:php]
autostart=true
command=php-fpm7 -F